# ubuntu_nifi_fargate_zk/playbooks/nodes/nifi.timer
[Unit]
Description=Starts NiFi after boot.

[Timer]
OnBootSec=1min
Unit=nifi.service

[Install]
WantedBy=multi-user.target
